# etc/cron.hourly/lightdm-dynamic-cron
#!/bin/bash

echo "[LightDM-dynamic] : Cron START"

/bin/bash /usr/bin/lightdm-dynamic

echo "[LightDM-dynamic] : Cron END"
